# scripts/launchagents/com.oncologyit.newsletter-generate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.oncologyit.newsletter-generate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jarvis2amrahs/projects/oncologyit/scripts/generate-newsletter.sh</string>
    </array>
    <!-- Every Monday at 5:00 AM CT -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jarvis2amrahs/logs/oncologyit-newsletter.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jarvis2amrahs/logs/oncologyit-newsletter-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/Users/jarvis2amrahs/.npm-global/bin:/usr/bin:/bin</string>
        <key>MC_ENDPOINT</key>
        <string>http://100.80.207.81:3777</string>
    </dict>
</dict>
</plist>
